# packages/vxscreenshots/vxscreenshots-2.7.2.tar.gz/vxscreenshots-2.7.2/vxscreenshots/config/autostart/vxssicon.desktop
[Desktop Entry]
Type=Application
Exec=vxssicon
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=VxSSIcon Autocreated
Name=VxSSIcon Autocreated
Comment[en_US]=VxSSIcon Autocreated
Comment=VxSSIcon Autocreated
