# packages/vxscreenshots/vxscreenshots-2.7.2.tar.gz/vxscreenshots-2.7.2/vxscreenshots/config/autostart/vxsswatcher.desktop
[Desktop Entry]
Type=Application
Exec=vxsswatcher
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=VXSSwatcher Autocreated
Name=VXSSwatcher Autocreated
Comment[en_US]=VxSSIconwatcher Autocreated
Comment=VxSSwatcher Autocreated
