# packages/easysched/easysched-1.2.9-py3-none-any.whl/easysched-1.2.9.data/data/share/easysched/session_simple/easysched.service
[Unit]
Description=easysched startup script

[Service]
ExecStart=esdd --config=share/example.cfg --log-level=info
ExecStop=esdctrl system exit
Type=simple

[Install]
WantedBy=multi-user.target
