# packages/ch_frb_wview/ch_frb_wview-0.3.9.tar.gz/ch_frb_wview-0.3.9/ch_frb_wview.service
[Unit]
Description=CHIME FRB Operations Web Application
After=network-online.target

[Service]
ExecStart=/usr/local/bin/uwsgi \
        --http :5002 \
        --manage-script-name \
        --mount /=ch_frb_wview:webapp \
        --threads 4
Restart=on-failure
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all
