# packages/easysched/easysched-1.2.9-py3-none-any.whl/easysched-1.2.9.data/data/share/easysched/session_xserver/easysched.service
[Unit]
Description=easysched startup script
PartOf=graphical-session.target

[Service]
ExecStart=esdd --config=share/example.cfg --log-level=INFO
ExecStop=esdctrl system exit
Type=simple
Environment="DISPLAY=localhost:0"

[Install]
WantedBy=graphical.target
